# Configs/udev/81-wol.rules
ACTION=="add", SUBSYSTEM=="net", NAME=="eth*", RUN+="/usr/sbin/ethtool -s $name wol g"
